# config/systemd/nvcommd.service
[Unit]
Description=NV Communication Daemon for RK3588
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/nvcommd -c /etc/nvcomm/nvcomm.conf
Restart=on-failure
WatchdogSec=30
CPUAffinity=0

[Install]
WantedBy=multi-user.target
